# old/old_prod_config/system/ddclient.service
[Unit]
Description=Dynamic DNS Update Client

[Service]
Type=oneshot
ExecStart=/usr/sbin/ddclient --force

[Install]
WantedBy=multi-user.target

